# Tweak/layout/DEBIAN/postinst
#!/bin/sh
launchctl load /Library/LaunchDaemons/org.thebigboss.cpconnector.plist
